# debian/prerm
#!/bin/bash

set -e

systemctl stop webzfs.service || true
systemctl disable webzfs.service || true
